# rootfs-overlay/usr/lib/systemd/system/rpc_pipefs.target
[Unit]
Requires=run-nfs-rpc_pipefs.mount
After=run-nfs-rpc_pipefs.mount
